# systemd/telegram-bot-failure.service
[Unit]
Description=Telegram Bot Crash Handler

[Service]
Type=oneshot
ExecStart=/bin/echo "telegram-bot crashed"
