# systemd/telegram-bot-update.service
[Unit]
Description=Update Telegram Bot from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/telegram-bot

ExecStart=/opt/telegram-bot/systemd/update.sh

Environment=PYTHONUNBUFFERED=1
TimeoutStartSec=300

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
