# systemd/telegram-bot-update.service
[Unit]
Description=Update Telegram Bot from GitHub

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /opt/telegram-bot/systemd/update.sh
