# systemd/telegram-bot-update.timer
[Unit]
Description=Auto-check Telegram Bot updates every 60 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min
Persistent=true

[Install]
WantedBy=timers.target
