# systemd/telegram-bot-update.service
[Unit]
Description=Telegram Bot GitHub Updater

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-bot
ExecStart=/home/ubuntu/telegram-bot/updater.sh
